# setup/local.ollama.system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.personalagent.ollama</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/local/bin/start_ollama.sh</string>
    </array>

    <!-- Ensure launchd has the env it needs -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>

      <!-- Tune as desired -->
      <key>OLLAMA_HOST</key>
      <string>0.0.0.0:11434</string>
      <key>OLLAMA_ORIGINS</key>
      <string>*</string>
      <key>OLLAMA_MAX_LOADED_MODELS</key>
      <string>2</string>
      <key>OLLAMA_NUM_PARALLEL</key>
      <string>8</string>
      <key>OLLAMA_MAX_QUEUE</key>
      <string>512</string>
      <key>OLLAMA_KEEP_ALIVE</key>
      <string>60m</string>
      <key>OLLAMA_NEW_ENGINE</key>
      <string>true</string>
      <key>OLLAMA_NEW_ESTIMATES</key><string>true</string>
      <key>OLLAMA_MODELS</key><string>/Users/Shared/personal_agent_data/ollama_models</string>

      <!-- Optional:
      <key>OLLAMA_LOG_LEVEL</key><string>INFO</string>
      <key>OLLAMA_MODELS</key><string>/Volumes/Fast/ollama-models</string>
      -->
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Keep the service up if it crashes/exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Simple log files (use log stream for detailed logs) -->
    <key>StandardOutPath</key>
    <string>/tmp/ollama.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ollama.err</string>
  </dict>
</plist>
